# systemd/tarot.path.template
[Unit]
Description=<?=PROGNAME?> state file watcher

[Path]
PathExists=<?=TRIGGER_FILE?>


[Install]
WantedBy=multi-user.target
